# build/webserver/newznab.sh.conf
[supervisord]
nodaemon=true

[program:newznab]
command=/newznab.sh
user=%(ENV_APACHE_RUN_USER)s
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
tderr_logfile=/var/log/newznab/newznab.sh.err.log
stdout_logfile=/var/log/newznab/newznab.sh.out.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:svn-download]
autostart=true
autorestart=false
command=/bin/bash -c "/newznab-setup.sh"
stdout_logfile=/var/log/newznab/newznab-setup.sh.out.log
